# supervisord-postgresql.conf
[program:postgres]
command=/usr/local/bin/supervisord-postgresql.sh
numprocs=1
autostart=true
autorestart=true
