# supervisord-apache2.conf
[program:apache2]
command=/usr/local/bin/supervisord-apache2.sh
numprocs=1
autostart=true
autorestart=true
