# supervisord-postgresql.conf
[program:postgres]
command=/supervisord-postgresql.sh
numprocs=1
autostart=true
autorestart=true